# backend/supervisor.conf
[supervisord]
nodaemon=true

[program:laravel]
command=php artisan serve --host=0.0.0.0 --port=9000
directory=/var/www/html
autostart=true
autorestart=true
user=www-data

